# system-monitor.service
[Unit]
Description=System Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tyson
WorkingDirectory=/home/tyson/Tools/monitoring-tools
ExecStart=/usr/bin/python3 /home/tyson/Tools/monitoring-tools/monitor.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/tyson/Tools/monitoring-tools/logs

[Install]
WantedBy=multi-user.target
