# vps-monitor.service
[Unit]
Description=VPS Server Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tyson
WorkingDirectory=/home/tyson/AiTools/monitoring-tools
ExecStart=/usr/bin/python3 /home/tyson/AiTools/monitoring-tools/vps_monitor.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/tyson/AiTools/monitoring-tools/logs

[Install]
WantedBy=multi-user.target
